# etc/udev/rules.d/99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/ravil/.local/bin/monitor_setup.sh"
